# systemd/mixxx.service
[Unit]
Description=Mixxx DJ software
PartOf=mixxx.target
Conflicts=maschinepi.target
After=sound.target openbox.service mixxx-audio-profile.service
Requires=openbox.service mixxx-audio-profile.service

[Service]
Type=simple
User=mpi
Environment=DISPLAY=:99
Environment=HOME=/home/mpi
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PIPEWIRE_LATENCY=256/48000
Environment=XCURSOR_THEME=DMZ-White
Environment=XCURSOR_SIZE=24
ExecStartPre=/usr/local/sbin/mpi-configure-mixxx-controller
ExecStartPre=-/bin/bash -c 'for i in $(seq 1 30); do wpctl status 2>/dev/null | grep -q "Maschine MK3" && exit 0; sleep 1; done; echo "MK3 audio not found after 30s; starting anyway"'
ExecStart=/usr/bin/pw-jack /usr/local/bin/mixxx --fullScreen
Restart=on-failure
RestartSec=5

[Install]
WantedBy=mixxx.target
