# systemd/mk3-mode-selector.service
[Unit]
Description=Select stored MK3 operating mode
Documentation=file:/opt/mpi-station/mode-selector/README.md
After=mpi-station-first-boot.service multi-user.target
ConditionPathExists=/var/lib/mpi-station/provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mpi-mode-select
