# systemd/mk3-t9-daemon.service
[Unit]
Description=Mixxx MK3 T9 text input daemon
PartOf=mixxx.target
Conflicts=maschinepi.target
After=xvfb.service mixxx.service
Requires=xvfb.service

[Service]
Type=simple
User=mpi
Environment=DISPLAY=:99
ExecStart=/usr/bin/python3 /opt/mpi-station/external/mixxx-mk3/pi-setup/mk3-t9-daemon.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=mixxx.target
