# systemd/mk3-mode-selector.service
[Unit]
Description=Select MK3 operating mode at boot
Documentation=file:/opt/mpi-station/mode-selector/README.md
After=multi-user.target
ConditionPathExists=/var/lib/mpi-station/provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mk3-mode-selector
