# systemd/mpi-prepare-data.service
[Unit]
Description=Divide SD-card remainder between Mixxx and MusicPI data
DefaultDependencies=no
After=systemd-remount-fs.service
Before=home-mpi-Music.mount home-mpi-maschinepi-samples.mount local-fs.target
ConditionPathExists=!/var/lib/mpi-station/data-layout-ready

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mpi-prepare-data-partitions
RemainAfterExit=yes
TimeoutStartSec=5min

[Install]
WantedBy=local-fs.target
