# systemd/mpi-station-first-boot.service
[Unit]
Description=Provision mpi-station application stacks
Documentation=file:/opt/mpi-station/image/README.md
Wants=network-online.target
After=network-online.target
Before=mk3-mode-selector.service
ConditionPathExists=!/var/lib/mpi-station/provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mpi-station-first-boot
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
